# SaintBot/Private EXE Loader/loader_sigma.yml
title: New Private EXE Loader
id: f47ac10b-58cc-4372-a567-0e02b2c3d479
status: test
description: |
    Detects the new AV/EDR/XDR bypass exe loader.
references:
    - https://github.com/infinitumitlabs/
    - https://infinitumit.com.tr
    - https://www.linkedin.com/company/infinitumlabs
author: Aziz Kaplan
date: 2024/05/06
modified: 2024/05/06
tags:
    - attack.t1566
    - attack.t1059
    - attack.t1204
    - attack.t1543
    - attack.t1027
    - attack.t1112
    - attack.t1102
    - attack.t1071.001
    - attack.t1489
logsource:
    categories: file_access,network_activity,ps_script
    product: windows
detection:
    selection_IP:
        DestinationIp:
            - '146.70.56.165'
        DestinationHostname: 'kurd.computer'
    selection_directory:
        
        FileName|contains:
            - '\AppData\Local\Temp\'
        FileName|endswith:
            - '.exe'
    selection_registry_1:
        TargetObject|contains: 
            - '\CurrentControlSet\Services\Tcpip\Parameters'
    selection_registry_2:
        TargetObject|contains: 
            - '\Windows\CurrentVersion\Internet Settings\Connections'
    condition: 1 of selection_IP or (selection_directory and selection_registry_2 and selection_registry_1)
falsepositives:
    - Unlikely
level: high